# systemd/daily_reboot.service
[Unit]
Description=Reboot the system (scheduled)

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl reboot